# services/weatherService/deploy/systemd/weather-scheduler.service
[Unit]
Description=Weather Data Collector Scheduler
Documentation=https://github.com/JokerTrickster/joker_backend
After=network.target mysql.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=weather
Group=weather
WorkingDirectory=/opt/weather-scheduler

# Environment configuration
EnvironmentFile=/opt/weather-scheduler/.env

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/weather

# Process management
ExecStart=/opt/weather-scheduler/bin/scheduler
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-scheduler

[Install]
WantedBy=multi-user.target
